# voice-assistant.service
[Unit]
Description=Voice Assistant (Pipecat)
After=network.target

[Service]
Type=simple
User=zammitjames
WorkingDirectory=/home/zammitjames/actions-runner/_work/app-voice-assistant/app-voice-assistant
Environment="PATH=/home/zammitjames/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/zammitjames/.local/bin/uv run bot.py
Restart=always
RestartSec=5
TimeoutStopSec=5
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
